# config/systemd/probe-sniffer.service
[Unit]
Description=WiFi Probe Request Sniffer
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/probe-sniffer
EnvironmentFile=/opt/probe-sniffer/.env
Environment="PATH=/opt/probe-sniffer/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Run setup script to configure monitor mode, then start sniffer
ExecStart=/opt/probe-sniffer/start.sh

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
